# om-env-kako/templates/default/kako.service.erb
[Unit]
Description=Kako Service
After=syslog.target

[Service]
Type=simple
Environment=LC_ALL=C.UTF-8
Environment=LANG=C.UTF-8
User=<%= @user %>
WorkingDirectory=<%= @dir %>
ExecStart=<%= @script %>
StandardOutput=syslog
StandardError=syslog

[Install]
WantedBy=multi-user.target
